# etc/profile.d/reform-kwin.sh
# TODO: The bug for which the below is a workaround is yet to be filed with upstream.
#       See http://mntre.com/reform-irc-logs/2025-09-15.log.html#t17:48:39
KWIN_FORCE_SW_CURSOR=1
